# pqdbus.service
[Unit]
Description=PQ wrapper DBUS service

[Service]
Type=dbus
BusName=io.FuriOS.PQ
ExecStart=/usr/libexec/pqdbus
